# mcp-inspector/rootfs/etc/services.d/mcp-inspector/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Cleanup script for MCP Inspector service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

bashio::log.info "MCP Inspector service stopped." 